# systemd/user/unclutter.service
[Unit]
Description=Unclutter screen

[Service]
ExecStart=/usr/bin/unclutter --timeout 1 --fork

[Install]
WantedBy=default.target

